# deploy/zakupomat.service
[Unit]
Description=Zakupomat Backend API
After=network.target mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/zakupomat/backend
Environment="PATH=/var/www/zakupomat/backend/venv/bin"
EnvironmentFile=/var/www/zakupomat/backend/.env
ExecStart=/var/www/zakupomat/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
